# docker/app/rootfs/etc/service/nginx/run
#!/bin/sh

exec 2>&1
source /env

mkdir -p /run/nginx
exec /usr/sbin/nginx -c '/etc/nginx/nginx.conf'
